# callhome.timer
[Unit]
Description=HTTP request to callhome PHP script to publish router IP address

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
